# deploy/supervisor.ini
[program:everimg]
directory=/path/to/everimg/
command=/path/to/everimg/bin/everimg.phar
numprocs=1
autostart=true
autorestart=unexpected
startsecs=10
stopwaitsecs=10
exitcodes=0,2
stopsignal=TERM
redirect_stderr=false
stdout_logfile=/path/to/everimg/var/log/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_events_enabled=false
stderr_logfile=/path/to/everimg/var/log/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_events_enabled=false
environment=CONF_FILE=./conf/release.ini,LOG_LEVEL=debug,LOG_TO_FILE=false